# data-backup.service
[Unit]
Description=Data backup

[Service]
Type=oneshot
ExecStart=/bin/bash /home/deck/.local/backup.sh
