# data-backup.timer
[Unit]
Description=Backup Steam Deck

[Timer]
OnBootSec=120
OnCalendar=00/12:15
Persistent=true

[Install]
WantedBy=timers.target
