# data-backup.timer
[Unit]
Description=Backup Steam Deck

[Timer]
OnCalendar=*-*-* 11:00:00
Persistent=true

[Install]
WantedBy=timers.target
